# ansible/files/puma.service
[Unit]
Description=Puma HTTP Server
After=network.target

[Service]
Type=simple
EnvironmentFile=/home/gcp_appUser/db_config
User=gcp_appUser
WorkingDirectory=/home/gcp_appUser/reddit
ExecStart=/bin/bash -lc 'puma'
Restart=always

[Install]
WantedBy=multi-user.target